# meta-phosphor/common/recipes-phosphor/obmc-phosphor-flash/obmc-mgr-download/org.openbmc.managers.Download.service
[Unit]
Description=Phosphor OpenBMC DBus service management daemon

[Service]
ExecStart=@SBINDIR@/download_manager.py
Environment="PYTHONUNBUFFERED=1"
