# meta-phosphor/common/recipes-phosphor/obmc-phosphor-flash/obmc-flash-bios/org.openbmc.control.Flash.service
[Unit]
Description=Phosphor OpenBMC BMC code update

[Service]
ExecStart=@SBINDIR@/bmc_update.py
Environment="PYTHONUNBUFFERED=1"
