# meta-phosphor/common/recipes-phosphor/obmc-phosphor-flash/obmc-flash-bmc/org.openbmc.control.BmcFlash.service
[Unit]
Description=Phosphor OpenBMC BIOS code update

[Service]
ExecStart=@SBINDIR@/flash_bios.exe
